# backup.timer
[Unit]
Description=backup-timer

[Timer]
OnUnitActiveSec=60s
OnBootSec=10s

[Install]
WantedBy=backup.target